# sources/meta-aurora/recipes-gui/cubegui/files/aurora.service.in
[Unit]
Description=Aurora GUI App
After=systemd-user-sessions.service

[Service]
ExecStart=%sbindir%/cubeManager -f /usr/bin/config.json
StandardOutput=syslog

[Install]
Alias=display-manager.service
